# signatures.json
{
  "version": "1.0.0",
  "description": "evtx-to-csv enrichment signature library. Edit and save, then click Refresh in the app.",
  "patterns": [
    {
      "name": "ps_encoded_command",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1027",
      "description": "PowerShell launched with -EncodedCommand/-enc: Base64-encodes the payload to hide the actual command from casual inspection and bypass script-level logging.",
      "regex": "(?i)\\b(?:powershell|pwsh)(?:\\.exe)?\\b.*(?:-enc\\b|-encodedcommand\\b)"
    },
    {
      "name": "ps_invoke_expression",
      "risk": "High",
      "tactic": "Execution",
      "mitre_id": "T1059.001",
      "description": "Invoke-Expression (IEX) executes an arbitrary string as code at runtime. Hallmark of fileless malware.",
      "regex": "(?i)\\b(?:iex|invoke-expression)\\s*[\\(\\$]"
    },
    {
      "name": "ps_download_cradle",
      "risk": "High",
      "tactic": "Command & Control",
      "mitre_id": "T1105",
      "description": "Remote script download via WebClient/DownloadString or Invoke-WebRequest. Pulls attacker code directly into memory.",
      "regex": "(?i)\\b(?:downloadstring|downloadfile|net\\.webclient|invoke-webrequest|iwr|curl|wget)\\b"
    },
    {
      "name": "ps_bits_transfer",
      "risk": "Medium",
      "tactic": "Command & Control",
      "mitre_id": "T1197",
      "description": "BITS abused for file download/upload. BITS transfers survive reboots and run as a background service.",
      "regex": "(?i)\\b(?:start-bitstransfer|bitsadmin\\s+/transfer)\\b"
    },
    {
      "name": "ps_execution_bypass",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1059.001",
      "description": "PowerShell execution policy bypass or hidden-window flags. Used to suppress restrictions and avoid visible console windows.",
      "regex": "(?i)\\b(?:executionpolicy\\s+bypass|-ep\\s+bypass|-nop\\b|-noninteractive\\b|-windowstyle\\s+hidden|-w\\s+hidden)\\b"
    },
    {
      "name": "ps_reflection_load",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1620",
      "description": "Reflective assembly loading injects a .NET DLL directly into memory without touching disk.",
      "regex": "(?i)\\b(?:\\[reflection\\.assembly\\]::load|add-type\\s+-assembly|\\[system\\.reflection\\.assembly\\]::load)\\b"
    },
    {
      "name": "ps_remoting",
      "risk": "Medium",
      "tactic": "Lateral Movement",
      "mitre_id": "T1021.006",
      "description": "PowerShell remoting used with a -ComputerName target — indicates lateral movement attempt via WinRM.",
      "regex": "(?i)\\b(?:invoke-command|enter-pssession|new-pssession)\\b.*-computername\\b"
    },
    {
      "name": "ps_network_socket",
      "risk": "High",
      "tactic": "Command & Control",
      "mitre_id": "T1071",
      "description": "Direct use of .NET network socket classes. Indicates a custom C2 channel or data exfiltration routine.",
      "regex": "(?i)\\b(?:net\\.sockets|tcpclient|udpclient|system\\.net\\.sockets|net\\.socket)\\b"
    },
    {
      "name": "ps_amsi_bypass",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1562.001",
      "description": "AMSI bypass attempt detected. Attackers patch amsi.dll in memory to prevent real-time scanning.",
      "regex": "(?i)\\b(?:amsiutils|amsi\\.dll|disableamsi|amsicontext|amsiinit|amsiscancontent|amsibypass)\\b"
    },
    {
      "name": "ps_decompression",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1027.002",
      "description": "PowerShell using Gzip/Deflate decompression. Often used to hide large scripts or second-stage payloads.",
      "regex": "(?i)\\b(?:IO\\.Compression\\.(?:GzipStream|DeflateStream)|\\[System\\.IO\\.Compression\\.CompressionMode\\]::Decompress)\\b"
    },
    {
      "name": "ps_amsi_bypass_reflection",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1562.001",
      "description": "Advanced AMSI bypass using reflection to access AmsiUtils. Targets the internal .NET implementation of AMSI.",
      "regex": "(?i)\\b(?:\\[Ref\\]\\.Assembly\\.GetType\\('System\\.Management\\.Automation\\.AmsiUtils'\\))\\b"
    },
    {
      "name": "ps_marshal_memory_manipulation",
      "risk": "Critical",
      "tactic": "Execution",
      "mitre_id": "T1059.001",
      "description": "PowerShell using System.Runtime.InteropServices.Marshal for direct memory access. Indicates shellcode injection or memory exploitation.",
      "regex": "(?i)\\[system\\.runtime\\.interopservices\\.marshal\\]::(?:copy|read|write|allochglobal|ptrto|structuretoptr)"
    },
    {
      "name": "lolbin_certutil_download",
      "risk": "High",
      "tactic": "Command & Control",
      "mitre_id": "T1105",
      "description": "certutil.exe used to download a remote file — a well-known LOLBin abuse technique.",
      "regex": "(?i)\\bcertutil(?:\\.exe)?\\b.*(?:-urlcache|-split|-f\\s+http)"
    },
    {
      "name": "lolbin_mshta",
      "risk": "High",
      "tactic": "Execution",
      "mitre_id": "T1218.005",
      "description": "mshta.exe executing an HTA file or remote URL. MSHTA runs HTML Applications with full scripting privileges.",
      "regex": "(?i)\\bmshta(?:\\.exe)?\\b.*(?:http|vbscript:|javascript:|\\.hta\\b)"
    },
    {
      "name": "lolbin_regsvr32_squiblydoo",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218.010",
      "description": "Squiblydoo attack: regsvr32.exe used to execute a remote COM scriptlet. Bypasses AppLocker.",
      "regex": "(?i)\\bregsvr32(?:\\.exe)?\\b.*(?:/i\\s*:\\s*http|scrobj\\.dll)"
    },
    {
      "name": "lolbin_rundll32_minidump",
      "risk": "Critical",
      "tactic": "Credential Access",
      "mitre_id": "T1003.001",
      "description": "rundll32.exe invoking comsvcs.dll MiniDump — used to dump LSASS process memory.",
      "regex": "(?i)\\brundll32(?:\\.exe)?\\b.*(?:comsvcs\\.dll|minidump)"
    },
    {
      "name": "lolbin_wmic_exec",
      "risk": "Medium",
      "tactic": "Execution",
      "mitre_id": "T1047",
      "description": "WMIC used for process creation or remote execution without spawning cmd.exe directly.",
      "regex": "(?i)\\bwmic(?:\\.exe)?\\b.*(?:process\\s+call\\s+create|/node\\s*:)"
    },
    {
      "name": "lolbin_cscript_wscript",
      "risk": "High",
      "tactic": "Execution",
      "mitre_id": "T1059.005",
      "description": "CScript/WScript executing a script file. Frequently used to run .vbs/.js droppers.",
      "regex": "(?i)\\b(?:cscript|wscript)(?:\\.exe)?\\b.*\\.(vbs|js|jse|vbe|wsf)\\b"
    },
    {
      "name": "lolbin_installutil",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218.004",
      "description": "InstallUtil.exe abused to execute arbitrary .NET code. Bypasses application whitelisting.",
      "regex": "(?i)\\binstallutil(?:\\.exe)?\\b"
    },
    {
      "name": "lolbin_cmstp",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218.003",
      "description": "CMSTP.exe exploited to execute arbitrary DLLs or COM scriptlets. Bypasses UAC.",
      "regex": "(?i)\\bcmstp(?:\\.exe)?\\b"
    },
    {
      "name": "lolbin_msiexec_remote",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218.007",
      "description": "msiexec.exe used to install a package from a remote URL. Allows execution of arbitrary code.",
      "regex": "(?i)\\bmsiexec(?:\\.exe)?\\b.*(?:/i\\s+http|/package\\s+http|quiet)"
    },
    {
      "name": "lolbin_psexec",
      "risk": "High",
      "tactic": "Lateral Movement",
      "mitre_id": "T1569.002",
      "description": "PsExec used for remote command execution. Frequently used by ransomware operators.",
      "regex": "(?i)\\bpsexec(?:\\.exe|c\\.exe)?\\b"
    },
    {
      "name": "lolbin_pcalua_execution",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218",
      "description": "Program Compatibility Assistant (pcalua.exe) used to proxy execution.",
      "regex": "(?i)\\bpcalua(?:\\.exe)?\\b.*-a\\b"
    },
    {
      "name": "lolbin_file_extraction",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218",
      "description": "Native Windows tools (makecab, expand, extrac32) used for file extraction or cabinet creation.",
      "regex": "(?i)\\b(?:makecab|expand|extrac32)(?:\\.exe)?\\b"
    },
    {
      "name": "defender_tamper",
      "risk": "Critical",
      "tactic": "Defense Evasion",
      "mitre_id": "T1562.001",
      "description": "Windows Defender configuration modified. Adding exclusion paths or disabling real-time monitoring.",
      "regex": "(?i)\\b(?:set-mppreference|add-mppreference)\\b.*(?:disablebehaviormonitoring|disablerealtimemonitoring|exclusionpath|disableioavprotection|disableantispyware)"
    },
    {
      "name": "shadow_copy_delete",
      "risk": "Critical",
      "tactic": "Impact",
      "mitre_id": "T1490",
      "description": "Shadow copy/backup deletion detected. Executed immediately before ransomware encryption.",
      "regex": "(?i)\\b(?:vssadmin\\s+delete\\s+shadows|wbadmin\\s+delete\\s+catalog|wmic\\s+shadowcopy\\s+delete)\\b"
    },
    {
      "name": "event_log_clear_cmd",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1070.001",
      "description": "Windows event log cleared via wevtutil or PowerShell. Destroying forensic evidence.",
      "regex": "(?i)\\b(?:wevtutil\\s+cl\\b|wevtutil\\s+clear-log|clear-eventlog|remove-eventlog)\\b"
    },
    {
      "name": "bcdedit_recovery_disabled",
      "risk": "High",
      "tactic": "Impact",
      "mitre_id": "T1490",
      "description": "bcdedit used to disable Windows recovery. Prevents booting into recovery mode.",
      "regex": "(?i)\\bbcdedit(?:\\.exe)?\\b.*(?:recoveryenabled\\s+no|safeboot)"
    },
    {
      "name": "credential_dump_tool",
      "risk": "Critical",
      "tactic": "Credential Access",
      "mitre_id": "T1003",
      "description": "Known credential dumping tool detected: Mimikatz, WCE, fgdump, pwdump, gsecdump, or LSASS dump.",
      "regex": "(?i)\\b(?:mimikatz|sekurlsa|procdump\\b.*lsass|ntdsutil\\b|vaultcmd\\b|wce\\.exe|fgdump|pwdump|gsecdump|invoke-mimikatz|lsassy)\\b"
    },
    {
      "name": "credential_dump_lsass",
      "risk": "Critical",
      "tactic": "Credential Access",
      "mitre_id": "T1003.001",
      "description": "Direct LSASS process dump detected. Provides all logged-on user credential hashes.",
      "regex": "(?i)\\b(?:procdump(?:\\.exe)?)\\b.*\\blsass\\b|\\blsass\\b.*minidump"
    },
    {
      "name": "remote_access_tool",
      "risk": "High",
      "tactic": "Command & Control",
      "mitre_id": "T1219",
      "description": "Unauthorized remote access tool detected (AnyDesk, ngrok, netcat, chisel, Cobalt Strike).",
      "regex": "(?i)\\b(?:anydesk|ngrok|ncat\\.exe|netcat\\b|socat\\b|chisel\\b|ligolo\\b|meterpreter|cobalt.?strike|c2\\.exe|beacon\\.exe)\\b"
    },
    {
      "name": "wmi_persistence_eventing",
      "risk": "High",
      "tactic": "Persistence",
      "mitre_id": "T1546.003",
      "description": "WMI event subscription created. Used to trigger malicious code when system events occur.",
      "regex": "(?i)\\b(?:__EventFilter|__EventConsumer|FilterToConsumerBinding)\\b"
    },
    {
      "name": "persistence_run_key",
      "risk": "High",
      "tactic": "Persistence",
      "mitre_id": "T1547.001",
      "description": "Registry Run/RunOnce key modification. Program executes at every user login.",
      "regex": "(?i)\\breg\\s+add\\b.*\\\\(?:run|runonce)\\b"
    },
    {
      "name": "persistence_schtask_create",
      "risk": "Medium",
      "tactic": "Persistence",
      "mitre_id": "T1053.005",
      "description": "Scheduled task created via schtasks /create. Malware survives reboots and re-executes.",
      "regex": "(?i)\\bschtasks(?:\\.exe)?\\b.*(?:/create\\b|/sc\\b)"
    },
    {
      "name": "persistence_service_create",
      "risk": "Medium",
      "tactic": "Persistence",
      "mitre_id": "T1543.003",
      "description": "Windows service created or reconfigured. Malicious services for persistent SYSTEM-level execution.",
      "regex": "(?i)\\b(?:sc\\s+create\\b|sc\\s+config\\b|new-service\\b)\\b"
    },
    {
      "name": "persistence_winlogon_shell_userinit",
      "risk": "High",
      "tactic": "Persistence",
      "mitre_id": "T1547.004",
      "description": "Winlogon Shell or Userinit registry key modified. Provides reliable persistence at login.",
      "regex": "(?i)\\breg\\s+add\\b.*\\\\Winlogon\\b.*(?:Shell|Userinit)\\b"
    },
    {
      "name": "persistence_ifeo_debugger",
      "risk": "Critical",
      "tactic": "Persistence",
      "mitre_id": "T1546.012",
      "description": "Image File Execution Options (IFEO) debugger hijack. Runs malware whenever a targeted app launches.",
      "regex": "(?i)\\breg\\s+add\\b.*\\\\Image\\s+File\\s+Execution\\s+Options\\b.*\\\\Debugger\\b"
    },
    {
      "name": "persistence_appinit_dlls",
      "risk": "High",
      "tactic": "Persistence",
      "mitre_id": "T1546.010",
      "description": "Persistence via AppInit_DLLs registry key. Loads DLL into every process linking with User32.dll.",
      "regex": "(?i)\\breg\\s+add\\b.*\\\\AppInit_DLLs\\b"
    },
    {
      "name": "persistence_startup_folder",
      "risk": "Medium",
      "tactic": "Persistence",
      "mitre_id": "T1547.001",
      "description": "File copied to the Windows Startup folder to run automatically at login.",
      "regex": "(?i)\\bcopy\\b.*\\bMicrosoft\\\\Windows\\\\Start\\s+Menu\\\\Programs\\\\Startup\\b"
    },
    {
      "name": "user_account_creation",
      "risk": "Medium",
      "tactic": "Persistence",
      "mitre_id": "T1136.001",
      "description": "New local user account created via command line. Fallback access point.",
      "regex": "(?i)\\bnet\\s+user\\s+.*\\s+/add\\b"
    },
    {
      "name": "user_added_to_admin_group",
      "risk": "High",
      "tactic": "Privilege Escalation",
      "mitre_id": "T1078.003",
      "description": "User added to the local Administrators group. Grants full control over the system.",
      "regex": "(?i)\\bnet\\s+localgroup\\s+administrators\\s+.*\\s+/add\\b"
    },
    {
      "name": "firewall_disabled",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1562.004",
      "description": "Windows Firewall disabled via netsh. Allows unrestricted network traffic.",
      "regex": "(?i)\\bnetsh\\s+advfirewall\\s+set\\s+(?:allprofiles|currentprofile|domainprofile|publicprofile|privateprofile)\\s+state\\s+off\\b"
    },
    {
      "name": "rdp_enabled_registry",
      "risk": "High",
      "tactic": "Lateral Movement",
      "mitre_id": "T1021.001",
      "description": "Remote Desktop (RDP) enabled via registry. Facilitates remote interaction.",
      "regex": "(?i)\\breg\\s+add\\b.*fDenyTSConnections\\b.*/t\\s+REG_DWORD\\b.*/d\\s+0\\b"
    },
    {
      "name": "service_stop_security",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1562.001",
      "description": "Security-related service stopped via net stop. Targets Defender, SENSE, or Firewall.",
      "regex": "(?i)\\bnet\\s+stop\\s+(?:windefend|sense|mpssvc)\\b"
    },
    {
      "name": "recon_whoami",
      "risk": "Low",
      "tactic": "Discovery",
      "mitre_id": "T1033",
      "description": "whoami executed — confirms the compromised user account and its privilege level.",
      "regex": "(?i)\\bwhoami(?:\\.exe)?\\b"
    },
    {
      "name": "recon_net_user_domain",
      "risk": "Medium",
      "tactic": "Discovery",
      "mitre_id": "T1087.002",
      "description": "net user /domain or net group — domain user/group enumeration.",
      "regex": "(?i)\\bnet\\s+(?:user\\s+/domain|group\\s+/domain|localgroup\\b|accounts\\b)"
    },
    {
      "name": "recon_systeminfo",
      "risk": "Low",
      "tactic": "Discovery",
      "mitre_id": "T1082",
      "description": "systeminfo executed — gathers OS version, patch level, and hotfixes.",
      "regex": "(?i)\\bsysteminfo(?:\\.exe)?\\b"
    },
    {
      "name": "recon_network_discovery",
      "risk": "Low",
      "tactic": "Discovery",
      "mitre_id": "T1016",
      "description": "Network configuration discovery (arp, route, netstat, ipconfig).",
      "regex": "(?i)\\b(?:arp\\s+-a|route\\s+print|netstat\\b|ipconfig\\s*/all)\\b"
    },
    {
      "name": "recon_tasklist",
      "risk": "Low",
      "tactic": "Discovery",
      "mitre_id": "T1057",
      "description": "tasklist or Get-Process used to enumerate running processes.",
      "regex": "(?i)\\b(?:tasklist(?:\\.exe)?|get-process\\b|get-service\\b)\\b"
    },
    {
      "name": "recon_password_search",
      "risk": "Medium",
      "tactic": "Credential Access",
      "mitre_id": "T1552.001",
      "description": "Searching for passwords in files via findstr. Look for plaintext credentials.",
      "regex": "(?i)\\bfindstr\\b.*/s.*/i\\b.*(?:password|pwd|admin|login|cred)\\b"
    },
    {
      "name": "data_staging_archiving",
      "risk": "Medium",
      "tactic": "Collection",
      "mitre_id": "T1560.001",
      "description": "Data staging using archive tools (7z, rar, zip). Indicates preparation for exfiltration.",
      "regex": "(?i)\\b(?:7z|rar|zip|tar)\\b.*(?:\\s+a\\s+|\\s+m\\s+).*\\.(?:7z|rar|zip|tar)\\b"
    },
    {
      "name": "heuristic_base64_in_command",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1027",
      "description": "Long Base64 string detected in a command field. Suggests obfuscated payload delivery.",
      "regex": "(?:[A-Za-z0-9+/]{40,}={0,2})"
    },
    {
      "name": "heuristic_temp_path_execution",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1036",
      "description": "Executable running from Temp, AppData, or Downloads path.",
      "regex": "(?i)(?:\\\\(?:temp|tmp|appdata\\\\(?:local|roaming)|downloads)\\\\[^\\s]+\\.(?:exe|dll|bat|ps1|vbs|js))\\b"
    },
    {
      "name": "heuristic_double_extension",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1036.007",
      "description": "File with double extension detected (e.g. invoice.pdf.exe). Masquerading technique.",
      "regex": "(?i)\\.\\w{2,4}\\.(exe|dll|bat|cmd|ps1|vbs|js|hta)\\b"
    },
    {
      "name": "lolbin_cmd_chaining",
      "risk": "Medium",
      "tactic": "Execution",
      "mitre_id": "T1059.003",
      "description": "cmd.exe /c or /r used for command chaining or proxy execution.",
      "regex": "(?i)\\bcmd(?:\\.exe)?\\b.*(?:/c|/r)\\b"
    },
    {
      "name": "lolbin_dotnet_proxies",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218",
      "description": "RegAsm or RegSvcs used to proxy execution of arbitrary .NET code.",
      "regex": "(?i)\\b(?:regasm|regsvcs)(?:\\.exe)?\\b"
    },
    {
      "name": "lolbin_odbcconf",
      "risk": "High",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218.008",
      "description": "odbcconf.exe used to execute a DLL or script. Less common LOLBin.",
      "regex": "(?i)\\bodbcconf(?:\\.exe)?\\b"
    },
    {
      "name": "lolbin_xwizard",
      "risk": "Medium",
      "tactic": "Defense Evasion",
      "mitre_id": "T1218",
      "description": "xwizard.exe used for arbitrary DLL loading to bypass security controls.",
      "regex": "(?i)\\bxwizard(?:\\.exe)?\\b"
    }
  ]
}
